# docker/workspace/crontab/schedule
* * * * * root php ${CONTAINER_CODE_PATH}/bin/artisan schedule:run >> /dev/null 2>&1